# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticmq]
command=/elasticmq/run

[program:saboteur]
command=/usr/bin/saboteur-agent
